# pkg/etc/systemd/system/dockflow.service
# /etc/systemd/system/dockflow.service
[Unit]
Description=DockFlow Daemon
Documentation=https://github.com/yourname/dockflow
After=docker.service
Requires=docker.service

[Service]
Type=simple
ExecStart=/usr/bin/dockflow daemon
Restart=always
RestartSec=5

# 安全 & 运行环境
User=root
WorkingDirectory=/var/lib/dockflow

# 日志
StandardOutput=journal
StandardError=journal

# 资源限制（可选）
LimitNOFILE=1048576

[Install]
WantedBy=multi-user.target
